# docker/service/build/crontab
* * * * * cd /var/www/html && /usr/local/bin/php artisan url:cron >> /var/www/html/storage/logs/laravel.txt
